# frontend-server/services/autoglm-post.service
[Unit]
Description=AutoGLM POST Link Server
After=network.target

[Service]
Type=simple
User=autoglm
WorkingDirectory=/opt/autoglm-frontend
EnvironmentFile=-/etc/autoglm/frontend.env
ExecStart=/usr/bin/env python3 /opt/autoglm-frontend/server_post.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
