# ocr-backend/services/autoglm-ocr-sync.service
[Unit]
Description=AutoGLM OCR Backend Service
After=network.target

[Service]
Type=simple
User=autoglm
WorkingDirectory=/opt/autoglm-ocr-backend
EnvironmentFile=-/etc/autoglm/ocr-backend.env
ExecStart=/usr/bin/env python3 /opt/autoglm-ocr-backend/ocr_backend_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
